# systemd/atlas-backup.service
[Unit]
Description=Atlas Backup Service
After=network.target atlas-api.service
Requires=atlas-api.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/dev/atlas
Environment=PYTHONPATH=/home/ubuntu/dev/atlas
Environment=ATLAS_DB_PATH=data/atlas.db
Environment=ATLAS_LOG_LEVEL=INFO

ExecStart=/usr/bin/python3 scripts/backup.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-backup

# Resource limits
MemoryMax=128M
CPUQuota=50%
TasksMax=5
LimitNOFILE=8192

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/dev/atlas/data /home/ubuntu/dev/atlas/backups

# Timeout for backup operations
TimeoutStopSec=300