# systemd/atlas-link-pipeline.timer
[Unit]
Description=Atlas Link Pipeline Timer (every 2 hours)

[Timer]
# Run every 2 hours
OnCalendar=*-*-* 00/2:00:00
# Run on boot after 5 minutes
OnBootSec=5min
# Persist timer across reboots
Persistent=true
# Add random delay up to 5 minutes to avoid thundering herd
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
