# systemd/atlas-url-fetcher.service
[Unit]
Description=Atlas Simple URL Fetcher
After=network.target

[Service]
Type=simple
User=khamel83
WorkingDirectory=/home/khamel83/github/atlas
Environment=PATH=/home/khamel83/github/atlas/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/khamel83/github/atlas/venv/bin/python scripts/simple_url_fetcher.py
Restart=always
RestartSec=60

# RESOURCE LIMITS - gentle pace, won't impact system
CPUQuota=15%
MemoryMax=512M
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-url-fetcher

[Install]
WantedBy=multi-user.target
