# systemd/atlas-progress.timer
[Unit]
Description=Atlas Progress Snapshot (hourly)

[Timer]
OnCalendar=*:00:00
Persistent=true

[Install]
WantedBy=timers.target
