# systemd/atlas-youtube.timer
[Unit]
Description=Run Atlas YouTube History Sync Daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
