# systemd/atlas-status-report.timer
[Unit]
Description=Update Atlas podcast status report hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
