# systemd/atlas-verify.service
[Unit]
Description=Atlas Content Verification
After=network.target

[Service]
Type=oneshot
User=khamel83
WorkingDirectory=/home/khamel83/github/atlas
ExecStart=/home/khamel83/github/atlas/venv/bin/python scripts/verify_content.py --report
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
